# repo2docker/buildpacks/autocommit.service
[Unit]
Description=Autocommit service

[Service]
Type=simple
ExecStart=/bin/bash /home/jovyan/autocommit.sh

[Install]
WantedBy=multi.user.target